# infra/systemd/quant-hft-daily-settlement.service
[Unit]
Description=Quant HFT Daily Settlement (user unit)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/
Environment=QUANT_HFT_ENV=prod
Environment=QUANT_ROOT=%h/workspace/quant_platform_HF
EnvironmentFile=-/etc/quant/platform.conf
ExecStart=/bin/bash -lc ': "${QUANT_ROOT:?QUANT_ROOT is required}"; TRADING_DAY=$(date +%F); "${QUANT_ROOT}/scripts/ops/run_daily_settlement.sh" --trading-day "${TRADING_DAY}" --execute'
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
