# deploy/systemd/quant-hft-data-pipeline.service
[Unit]
Description=Quant HFT Data Pipeline
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/
Environment=QUANT_ROOT=/opt/quant_platform_HF
Environment=PYTHON_BIN=/usr/bin/python3
Environment=ARCHIVE_ENDPOINT=localhost:9000
Environment=ARCHIVE_BUCKET=quant-archive
Environment=ARCHIVE_PREFIX=etl
EnvironmentFile=-/etc/quant/platform.conf
EnvironmentFile=-/etc/quant/quant-hft-data-pipeline.env
ExecStart=/bin/bash -lc ': "${QUANT_ROOT:?QUANT_ROOT is required}"; "${PYTHON_BIN}" "${QUANT_ROOT}/scripts/data_pipeline/run_pipeline.py" --analytics-db "${QUANT_ROOT}/runtime/analytics.duckdb" --export-dir "${QUANT_ROOT}/runtime/exports" --archive-endpoint "${ARCHIVE_ENDPOINT}" --archive-bucket "${ARCHIVE_BUCKET}" --archive-prefix "${ARCHIVE_PREFIX}" --archive-local-dir "${QUANT_ROOT}/runtime/archive" --iterations 0 --interval-seconds 60'
Restart=always
RestartSec=5
LimitNOFILE=65535
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
