# deploy/systemd/quant-hft-core-engine.service
[Unit]
Description=Quant HFT Core Engine (bootstrap)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/
Environment=QUANT_ROOT=/opt/quant_platform_HF
Environment=CTP_CONFIG_PATH=
EnvironmentFile=-/etc/quant/platform.conf
EnvironmentFile=-/etc/quant/quant-hft-core-engine.env
ExecStart=/bin/bash -lc ': "${QUANT_ROOT:?QUANT_ROOT is required}"; CONFIG_PATH="${CTP_CONFIG_PATH:-${QUANT_ROOT}/configs/sim/ctp.yaml}"; "${QUANT_ROOT}/build/core_engine" "${CONFIG_PATH}"'
Restart=always
RestartSec=2
LimitNOFILE=65535
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
